# deploy/tatatask.service
[Unit]
Description=TaTaTask Kanban Web App
After=network.target

[Service]
# .NET integrates with systemd via Host.UseSystemd() (Type=notify readiness + journal logging)
Type=notify
WorkingDirectory=/opt/tatatask
ExecStart=/opt/tatatask/TaTaTask
Restart=always
RestartSec=10
TimeoutStopSec=30
SyslogIdentifier=tatatask
User=tatatask
Group=tatatask
Environment=ASPNETCORE_ENVIRONMENT=Production
EnvironmentFile=-/opt/tatatask/ssl/pass.env
# 默认使用 appsettings.json 里的 Kestrel 端口配置 (HTTP:5000 + HTTPS:5001)
# 如需覆盖，取消下一行注释并修改
# Environment=ASPNETCORE_URLS=http://0.0.0.0:5000
# Logs go to the systemd journal (view: journalctl -u tatatask -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
